# roles/selfcheck/templates/vpn-selfcheck.service.j2
# {{ ansible_managed }}
[Unit]
Description=VPN gateway self-check
After=network-online.target sing-box.service nftables.service
Wants=sing-box.service
{% if selfcheck_onfailure_unit | length > 0 %}
OnFailure={{ selfcheck_onfailure_unit }}
{% endif %}

[Service]
Type=oneshot
ExecStart={{ selfcheck_script_path }}
# Вывод в journal — `journalctl -u vpn-selfcheck` показывает причину фейла,
# а notify-tg прикладывает tail журнала к сообщению. Файл-лог
# {{ selfcheck_log_path }} дублируется самим скриптом через tee.
StandardOutput=journal
StandardError=journal
Nice=10
